# packaging/scripts/postinstall.sh
systemctl daemon-reload
systemctl enable cloudwatch-metrics-proxy.service
systemctl restart cloudwatch-metrics-proxy.service
