# packaging/scripts/preremove.sh
systemctl stop cloudwatch-metrics-proxy.service
systemctl disable cloudwatch-metrics-proxy.service
